# private_dot_config/systemd/user/walker.service

[Unit]
Description=Walker

[Service]
ExecStart=/usr/bin/walker --gapplication-service
Restart=on-failure

[Install]
WantedBy=default.target
